# py/lora-gateway.conf
[program:lora-gateway]
directory=/home/pi/lora-pi-arduinos/py
command=python lora-gateway.py -d /dev/ttyACM0 -r 9600
autostart=true
autorestart=true
